# call-recorder.service
[Unit]
Description='Call Recorder'

[Service]
ExecStart=%h/bin/call-recorder
Restart=on-failure

[Install]
WantedBy=default.target
